# check-resolv.service
[Unit]
After=sshd.service

[Service]
ExecStart=/usr/local/bin/check-resolv.sh

[Install]
WantedBy=default.target
